# ytmusic_bot.service
[Unit]
Description=YouTube Music Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/usr/bin/python3 /home/ubuntu/ytmusic_interactive_bot.py
Restart=always
RestartSec=10

# Environment variables (optional)
# Environment="TOKEN=your_bot_token_here"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ytmusic-bot

[Install]
WantedBy=multi-user.target
